# supervisord.conf
[supervisord]
nodaemon=true

[program:worker]
command=./worker

[program:server]
command=./server
